# configs/baseline/airootfs/etc/systemd/system/etc-keys.automount
[Unit]
Description=Keys directory bind automount

[Automount]
Where=/etc/keys
TimeoutIdleSec=3min

[Install]
WantedBy=multi-user.target
