# configs/baseline/airootfs/etc/systemd/system/home.automount
[Unit]
Description=Home partition automount

[Automount]
Where=/home
TimeoutIdleSec=3min

[Install]
WantedBy=multi-user.target
