# configs/baseline/airootfs/etc/systemd/system/var-lib-systemd-home.automount
[Unit]
Description=systemd-homed credentials automount

[Automount]
Where=/var/lib/systemd/home
TimeoutIdleSec=3min

[Install]
WantedBy=multi-user.target
