# configs/baseline/airootfs/etc/systemd/system/var-lib-AccountService-users.automount
[Unit]
Description=users option automount

[Automount]
Where=/var/lib/AccountService/users
TimeoutIdleSec=3min

[Install]
WantedBy=multi-user.target
