# configs/baseline/airootfs/etc/systemd/system/var-lib-AccountsService-users.automount
[Unit]
Description=users option automount

[Automount]
Where=/var/lib/AccountsService/users
TimeoutIdleSec=3min

[Install]
WantedBy=multi-user.target
