# configs/baseline/airootfs/etc/systemd/system/var-lib-bluetooth.automount
[Unit]
Description=Bluetooth data automount

[Automount]
Where=/var/lib/bluetooth
TimeoutIdleSec=3min

[Install]
WantedBy=multi-user.target
